# ffda/gluon-channel-survey/files/etc/init.d/channel-survey
#!/bin/sh /etc/rc.common

START=99

boot() {
	# populate channel-survey on boot
	/usr/sbin/channel-survey > /tmp/gluon/channel-survey.cache
}
