# ffda/gluon-channel-survey/files/etc/init.d/channel-survey
#!/bin/sh /etc/rc.common

START=99

boot() {
	# populate channel-survey on boot, give radio some time to get ready...
	(sleep 30; /usr/sbin/channel-survey > /tmp/gluon/channel-survey.cache) &
}
